# deployments/systemd/oms-engine.service
[Unit]
Description=OMS C++ Core Engine
Documentation=https://github.com/mExOms/oms
After=network.target nats.service

[Service]
Type=simple
User=oms
Group=oms
WorkingDirectory=/opt/oms

# Binary and configuration
ExecStart=/opt/oms/bin/oms-engine
Environment="CONFIG_FILE=/opt/oms/configs/engine.yaml"

# Performance tuning
CPUAffinity=2 3
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Memory locking
LimitMEMLOCK=infinity
PrivateDevices=false
ProtectSystem=false

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oms-engine

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ReadWritePaths=/opt/oms/data /dev/shm

[Install]
WantedBy=multi-user.target